# deploy/clawbridge.service
[Unit]
Description=ClawBridge - Bridge your OpenClaw agents across machines
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/clawbridge
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=A2A_PORT=9100
Environment=A2A_BIND=0.0.0.0
EnvironmentFile=-/opt/clawbridge/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/clawbridge/config /opt/clawbridge/logs
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clawbridge

[Install]
WantedBy=multi-user.target
